# assistant/en-us.supervisord.conf
[supervisord]
nodaemon=true

[program:audio_input]
command=pulseaudio-input --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true

[program:wake_word]
command=porcupine-mqtt --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true

[program:voice_command]
command=webrtcvad-mqtt --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true

[program:speech_to_text]
command=pocketsphinx-mqtt --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true

[program:intent_recogntion]
command=fsticuffs-mqtt --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true

[program:text_to_speech]
command=espeak-mqtt --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true

[program:audio_output]
command=pulseaudio-output-mqtt --profile "%(ENV_PWD)s"
killasgroup=true
stopasgroup=true