# assistant/en-us.supervisord.conf
[supervisord]
nodaemon=true

[program:audio_input]
command=rhasspy-pulseaudio-input --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:wake_word]
command=rhasspy-porcupine-mqtt --profile "%(ENV_profile_dir)s" --auto-start --debug
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:voice_command]
command=rhasspy-webrtcvad-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:speech_to_text]
command=rhasspy-pocketsphinx-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:intent_recogntion]
command=rhasspy-fsticuffs-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:text_to_speech]
command=rhasspy-espeak-mqtt --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:audio_output]
command=rhasspy-pulseaudio-output-mqtt --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:training]
command=rhasspy-train-mqtt --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true