# assistant/en-us.supervisord.conf
[supervisord]
nodaemon=true

[program:audio_input]
command=rhasspy-pulseaudio-input --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/audio_input.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:wake_word]
command=rhasspy-porcupine-mqtt --profile "%(ENV_profile_dir)s" --auto-start --debug
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/wake_word.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:voice_command]
command=rhasspy-webrtcvad-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/voice_command.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:speech_to_text]
command=rhasspy-pocketsphinx-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/speech_to_text.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:intent_recogntion]
command=rhasspy-fsticuffs-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/intent_recognition.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:text_to_speech]
command=rhasspy-espeak-mqtt --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/text_to_speech.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:audio_output]
command=rhasspy-pulseaudio-output-mqtt --profile "%(ENV_profile_dir)s"
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/audio_output.log
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:training]
command=rhasspy-train-mqtt --profile "%(ENV_profile_dir)s" --debug
killasgroup=true
stopasgroup=true
stdout_logfile=%(ENV_log_dir)s/training.log
stdout_logfile_maxbytes=0
redirect_stderr=true